# pebble/kode-nas-display.service
[Unit]
Description=KODE NAS pebble OLED display daemon
Documentation=https://github.com/kode-nas/kode-os
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kode
Group=kode
# gpio + spi own the device nodes the daemon needs (gpiochip0, gpiomem0,
# spidev0.0). Listed here in addition to usermod -aG so the unit works
# even when group state hasn't been re-read for the kode user yet.
SupplementaryGroups=gpio spi
# Force the rpigpio backend. On a Pi 5 + Bookworm install, RPi.GPIO is
# actually python3-rpi-lgpio (a chardev-only shim) — it talks to
# /dev/gpiochip0 directly and doesn't need an lgd daemon. The plain
# `lgpio` factory looks like the obvious pick but the python3-lgpio
# package on Bookworm dies at import time trying to open an lgd
# notification FIFO that doesn't exist (no lgd daemon is shipped).
Environment=GPIOZERO_PIN_FACTORY=rpigpio
ExecStart=/usr/bin/python3 /home/kode/kode_nas_display.py
Restart=on-failure
RestartSec=5
# Hardening — the daemon only reads SPI + /proc + posts HTTP to local apps.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
DevicePolicy=closed
DeviceAllow=/dev/spidev0.0 rw
# Pi 5 RP1 GPIO chip + legacy /dev/gpiomem0 fallback for older lgpio.
DeviceAllow=/dev/gpiochip0 rw
DeviceAllow=/dev/gpiomem0 rw

[Install]
WantedBy=multi-user.target
